# wake-dryrun.service
[Unit]
Description=Wake-word dry-run logger (log-only, no execution path)

[Service]
ExecStart=/usr/lib/agent-fleet/wake-dryrun
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
